# ai.openclaw.memory-maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.memory-maintenance</string>
    
    <key>Comment</key>
    <string>OpenClaw nightly memory maintenance</string>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/sydneyassistent/bin/nightly-memory-maintenance.sh</string>
    </array>
    
    <key>StandardOutPath</key>
    <string>/Users/sydneyassistent/.openclaw/logs/memory-maintenance.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sydneyassistent/.openclaw/logs/memory-maintenance.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/sydneyassistent</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Users/sydneyassistent/.nvm/versions/node/v22.22.0/bin</string>
    </dict>
</dict>
</plist>